# scripts/crontab
#crontab file
#launches the script every ten minutes
SHELL=/bin/bash

*/10 * * * * cd /home/etud/Documents/GitHub/SAE23/scripts && bash script_lancement

